# gentoo-update.timer
[Unit]
Description=Gentoo System Update-Timer

[Timer]
OnCalendar=12:00
Persistent=true

[Install]
WantedBy=timers.target

