# gentoo-update.timer
[Unit]
Description=Gentoo System Update-Timer

[Timer]
OnCalendar=15:35
Persistent=true

[Install]
WantedBy=timers.target

